# etc/cron.d/zenloadbalancer
*/5 * * * *	root	/usr/local/zenloadbalancer/app/zenrrd/zenrrd.pl &>/dev/null
00  * * * * 	root	/usr/local/zenloadbalancer/app/zenntp/zenntp.pl &>/dev/null
